# deploy/truewebwork-watchdog.timer
[Unit]
Description=Run TrueWeb watchdog every minute

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=truewebwork-watchdog.service

[Install]
WantedBy=timers.target
